# bin/tst/org.apache.tomcat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><!-- NOTE: place this file in /Library/LaunchDaemons -->
    <dict>
        <key>Label</key>
        <string>org.spotweb.sync</string>
        <key>OnDemand</key>
        <false/>
        <key>RunAtLoad</key>
        <true/>
        <key>EnvironmentVariables</key>
        <dict>
            <key>WRAPPED_HOME</key>
            <string>/usr/local/bin</string>
        </dict>
        <key>ProgramArguments</key>
        <array>
            <string>/usr/local/bin/launchd_wrapper.sh</string>
        </array>
        <key>ServiceDescription</key>
        <string>Spotweb</string>
        <key>StandardErrorPath</key>
        <string>/tmp/spotweb.stderr</string>
        <key>StandardOutPath</key>
        <string>/tmp/spotweb.stdout</string>
    </dict>
</plist>
